# rss/rss.timer
[Unit]
Description=Run RSS Feed Parser every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
