# vdsalert/vdsalert.timer
[Unit]
Description=vdsalert timer

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
